# ansible/playbooks/templates/supervisor.ini
[program:{{ app_name }}]
command=/usr/local/bin/uwsgi
  --ini /etc/uwsgi/{{ app_name }}.ini
directory={{ webapps_dir }}/{{ app_name }}/src
autostart=true
autorestart=true
